# app/templates/service/_.service
[Service]
Type=simple
WorkingDirectory=/opt/dive/<%= package %>/
ExecStart=/opt/dive/<%= package %>/start.sh